# backup/echo.timer
[Unit]
Description=Run echo

[Timer]
OnCalendar=*:*:10
AccuracySec=1s

[Install]
WantedBy=timers.target
