# backup/echo.service
[Unit]
Description=borgmatic backup

[Service]
Type=simple
ExecStart=/home/abergman/.local/bin/hello-world.sh
