# backup/echo.service
[Unit]
Description=echo service

[Service]
Type=simple
ExecStart=/home/abergman/.local/bin/hello-world.sh
